# systemd/shelley-active-crawl.timer
[Unit]
Description=Daily Active Knowledge Crawl

[Timer]
OnCalendar=*-*-* 11:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
