# systemd/wiki-server.service
[Unit]
Description=Wiki Browser Server
After=network-online.target

[Service]
Type=simple
User=exedev
ExecStart=/home/exedev/.hermes/hermes-agent/venv/bin/python /home/exedev/scripts/wiki_server.py
Restart=always
RestartSec=5
Environment=PATH=/home/exedev/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
